# systemd/nodelistdb-testdaemon.service
[Unit]
Description=NodelistDB FidoNet Test Daemon
Documentation=https://github.com/xx25/nodelistdb
After=network.target
Wants=network.target

[Service]
Type=simple
User=daemon
Group=daemon
WorkingDirectory=/opt/nodelistdb
ExecStart=/opt/nodelistdb/testdaemon -config /etc/nodelistdb/testdaemon.yaml
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodelistdb-testdaemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/cache/nodetest /var/log/nodetest
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Resource limits
LimitNOFILE=65536
MemoryMax=4G

# Network access required for testing nodes
PrivateNetwork=false

[Install]
WantedBy=multi-user.target
